# scripts/com.bullcore.autorenew.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bullcore.autorenew</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/karlo.deltoro/Desktop/tasty-alerts/venv/bin/python</string>
        <string>/Users/karlo.deltoro/Desktop/tasty-alerts/auto_renew.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/karlo.deltoro/Desktop/tasty-alerts</string>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/karlo.deltoro/Desktop/tasty-alerts/auto_renew.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/karlo.deltoro/Desktop/tasty-alerts/auto_renew.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/karlo.deltoro/Desktop/tasty-alerts/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
